# internal/configs/systemd/ddns6-update.service
[Unit]
Description=DDNS6 Auto Update Service
After=network.target

[Service]
Type=simple
Environment=
ExecStart=test
WorkingDirectory=